# extras/supervisord-s.conf
[program:s]
command=/opt/release/bin/s
redirect_stderr=true
stdout_logfile=/var/log/s.log
